# supervisord.conf
[supervisord]
nodaemon=true

[program:webui]
command=node /webui/production_server.js
stdout_logfile = /tmp/webui.log
stderr_logfile = /tmp/webui-stderr.log
